# dist/system/kbs.service
[Unit]
Description=Key-Broker-Service
Documentation=https://confidentialcontainers.org
After=network.target
Requires=as.service

[Service]
ExecStart=/usr/local/bin/kbs --config-file /etc/trustee/kbs-config.toml
Environment=RUST_LOG=debug
Restart=always
RestartSec=5
Delegate=yes
KillMode=process
OOMScoreAdjust=-999
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity

[Install]
WantedBy=multi-user.target